# deploy/quadlet/templates/tidepool-proxy.container.in
# deploy/quadlet/templates/tidepool-proxy.container.in
# Rendered by deploy/scripts/render.sh — do not install the template itself.
#
# The frontend proxy: the ONLY network-facing process. Forwards TCP to the
# API's Unix socket with size/time limits; holds no data whatsoever — its
# mount set is run/ (for api.sock) and nothing else. Host exposure is the
# PublishPort line: the listen address is 127.0.0.1 unless the operator has
# deliberately configured a trusted-LAN address at render time.
#
# The proxy's own netns has outbound routes (rootless networking has no
# no-egress inbound-only mode) — the proxy CODE never dials anything but the
# fixed socket path — a code-level property of the proxy;
# the host. That layering — not the container flag — is the enforcement.

[Unit]
Description=Tidepool frontend proxy (local-network TCP → API Unix socket)
StartLimitIntervalSec=300
StartLimitBurst=5

[Container]
Image=@PROXY_IMAGE@
ContainerName=tidepool-proxy

ReadOnly=true
Tmpfs=/tmp:rw,noexec,nosuid,size=16m
Volume=@TIDEPOOL_HOME@/run:/var/lib/tidepool/run:rw
Volume=@TIDEPOOL_HOME@/config:/var/lib/tidepool/config:ro
# TLS: ONLY the server subtree (cert/key/chain), read-only. The CA private
# key under tls/ca/ is deliberately NOT mounted — see deploy/README.md.
Volume=@TIDEPOOL_HOME@/tls/server:/var/lib/tidepool/tls/server:ro

PublishPort=@LISTEN_ADDR@:@LISTEN_PORT@:@INTERNAL_PORT@

UserNS=keep-id:uid=@CONTAINER_UID@,gid=@CONTAINER_GID@
DropCapability=all
NoNewPrivileges=true
PodmanArgs=--pids-limit=64
PodmanArgs=--memory=128m --cpus=0.5 --ulimit=nofile=2048:2048

# in-container bind must be the netns wildcard for the published-port path;
# host reachability is confined by PublishPort above
Environment=TIDEPOOL_PROXY_ADDR=0.0.0.0
Environment=TIDEPOOL_PROXY_PORT=@INTERNAL_PORT@
Environment=TIDEPOOL_QUADLET_DIGEST=@QUADLET_DIGEST@

# health over the proxy's own liveness endpoint; NODE_TLS_REJECT is scoped
# to this one probe process (the local self-signed/local-CA cert) and does
# not affect request handling
HealthCmd=node -e "const https=require('https');https.get({host:'127.0.0.1',port:@INTERNAL_PORT@,path:'/health/live',rejectUnauthorized:false},r=>process.exit(r.statusCode===200?0:1)).on('error',()=>process.exit(1))"
HealthInterval=30s
HealthRetries=3
HealthStartPeriod=15s
HealthTimeout=5s

PodmanArgs=--stop-timeout=15

[Service]
Restart=on-failure
RestartSec=10
TimeoutStopSec=120

[Install]
WantedBy=default.target
